# deploy/macos/com.claudemem.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudemem.worker</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USER/Applications/claude-mem/ClaudeMem.Worker</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/USER/Applications/claude-mem</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_MEM_WORKER_PORT</key>
        <string>37777</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/claude-mem-worker.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/claude-mem-worker.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
